# roles/task_scheduler/templates/task.timer.j2
[Timer]
OnCalendar={{ task.schedule }}
Persistent=true

[Install]
WantedBy=timers.target